# config/supervisord.ini
[program:feedreader]
command=node /home/feedreader/apps/feedreader/src/server.js
numprocs=2
process_name=feedreader-%(process_num)s
autorestart=true
user=feedreader
stdout_logfile=/var/log/my-apps/feedreader/app.log
redirect_stderr=true
environment=PORT="2003%(process_num),ENV_PATH=/home/feedreader/config/feedreader.env"
directory=/home/feedreader/apps/feedreader
